# files/root.mount
[Unit]
Before=snapd.socket

[Mount]
What=/var/roothome
Where=/root
Options=bind
Type=none

[Install]
WantedBy=snapd.socket
